# bash/crontab
* * * * * /usr/local/bin/php -f /var/www/html/bin/cron/notifications.php
